# udev/rules.d/99-icpc.rules
SUBSYSTEM=="usb", DRIVER=="usb-storage", RUN+="/etc/icpc/udev/icpc-usb-check"
